# rootfs/etc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

# ── PostgreSQL ─────────────────────────────────────────────────────────────
[program:postgresql]
priority=100
command=su-exec postgres postgres -D /data/postgresql
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=QUIT
stopwaitsecs=30

# ── Database migrations (runs once after PostgreSQL is ready) ──────────────
[program:migrate]
priority=200
command=/usr/local/bin/run-migrate.sh
autostart=true
autorestart=false
startsecs=0
startretries=5
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── PHP-FPM ────────────────────────────────────────────────────────────────
[program:php-fpm]
priority=300
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Nginx ──────────────────────────────────────────────────────────────────
[program:nginx]
priority=400
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Symfony Messenger worker ───────────────────────────────────────────────
[program:worker]
priority=300
command=php /var/www/html/bin/console messenger:consume async --time-limit=3600 -vv
autostart=true
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
